# squad/pl140/systemd/line-watchdog-heartbeat.timer
[Unit]
Description=X63 run watchdog heartbeat watcher every 5 minutes (independent of the watchdog itself)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=line-watchdog-heartbeat.service

[Install]
WantedBy=timers.target
